# ansible/roles/spring-boot/templates/balruno-backend.service.j2
# {{ ansible_managed }}
# Balruno Spring Boot backend systemd unit.

[Unit]
Description=Balruno Spring Boot backend (Java 25)
After=network-online.target
Wants=network-online.target
# Don't fail to bring up the whole boot sequence if the JAR isn't deployed yet —
# the unit just stays inactive until the first backend-deploy run.
ConditionPathExists={{ sb_jar_path }}

[Service]
Type=simple
User={{ sb_user }}
Group={{ sb_user }}
WorkingDirectory={{ sb_install_dir }}

# JVM options + Spring profile passed as JVM args.
Environment="JAVA_TOOL_OPTIONS={{ sb_jvm_opts }}"
Environment="SPRING_PROFILES_ACTIVE={{ sb_active_profile }}"

ExecStart=/usr/bin/java -jar {{ sb_jar_path }}

# Logs go to journalctl by default; rotation handled by journald.
StandardOutput=journal
StandardError=journal
SyslogIdentifier={{ sb_service_name }}

Restart=on-failure
RestartSec=5
SuccessExitStatus=143

# Hardening
ProtectSystem=full
ProtectHome=true
NoNewPrivileges=true
PrivateTmp=true
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
